# gltail.crontab
# min hour dom month dow command
  0   */3  *   *     *   /Users/bheeshmar/gltail/restart_bubbles.sh
