# gltail.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
   <key>Label</key>
   <string>com.user.bubbles</string>
   <key>Program</key>
   <string>/Users/bheeshmar/gltail/restart_bubbles.sh</string>
   <key>RunAtLoad</key>
   <true/>
   <key>StartCalendarInterval</key>
   <array>
      <dict>
         <key>Hour</key>
         <integer>0</integer>
      </dict>
      <dict>
         <key>Hour</key>
         <integer>12</integer>
      </dict>
   </array>
</dict>
</plist>
